# BSP/src/multimedia/mdpd3/v4l2_mdpd/lib/systemd/system/mdpd.service
[Unit]
Description=MDP Daemon
DefaultDependencies=no
After=systemd-udevd.service
Before=sysinit.target

[Service]
ExecStart=/usr/bin/mdpd -f
Restart=always
StandardOutput=tty
StandardError=tty

[Install]
WantedBy=sysinit.target
