# BSP/meta/meta-mediatek-mt8110/recipes-support/nvram/files/nvram_daemon.service
[Unit]
Description=nvram_daemon

[Service]
ExecStart=/sbin/nvram_daemon
Type=simple

[Install]
Alias=nvramdaemon
WantedBy=multi-user.target
