# BSP/src/connectivity/combo_tool/mt7668_drv_insmod.service
[Unit]
Description=mt7668_drv_insmod


[Service]
Type=simple
ExecStart=/bin/sh /etc/mt7668_drv_insmod.sh

[Install]
WantedBy=multi-user.target
